# scripts/supervisor-fisherman.conf
[program:fisherman]
command=bash /opt/fisherman/run.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/fisherman.err.log
stdout_logfile=/var/log/fisherman.out.log
user=root
stopasgroup=true